# docs/firmware/pi_robot/pi-robot-rtc.service
[Unit]
Description=Better Robotics — Pi WebRTC peer (aiortc)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Runs as the robot user — the PTY spawned for `shell` channels lands at
# this UID, matching `ssh robot@<host>.local`. Independent of pi-robot
# (which still runs as root for BLE), so a crash here doesn't touch
# motors / LEDs.
User=robot
# Phase 2.F.1: BLE-signaled WebRTC offers arrive over the Unix socket at
# /run/pi-robot-rtc/sock, bridged here from pi_robot.py (which owns the
# BLE GATT). The robot user can't write to /run/ directly, so let
# systemd create the per-service runtime directory as robot:robot — the
# socket lives inside, 0o666 set by the daemon so root (pi_robot.py)
# can connect from outside.
RuntimeDirectory=pi-robot-rtc
RuntimeDirectoryMode=0755
WorkingDirectory=__HOME__/better-robotics/firmware/pi_robot
ExecStart=__HOME__/better-robotics/firmware/pi_robot/.venv/bin/python pi_robot_rtc.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
